# file/.config/autostart/xcompmgr (composer).desktop
[Desktop Entry]
Exec=xcompmgr &
Name=xcompmgr (composer)
Type=Application
Version=1.0
